# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#[ -f /usr/bin/btnd ] && /usr/bin/btnd wps 2 &
[ -f /usr/bin/btnd ] && /usr/bin/btnd reset 1 &
[ -f /usr/bin/watchdog ] && /usr/bin/watchdog &

[ -x /usr/sbin/dropbear ] || passwd -d root
[ -x /usr/sbin/dropbear ] || chmod +x /etc/config/telnetd
[ -x /usr/sbin/dropbear ] || /etc/config/telnetd -p 7777

exit 0

